# .config/systemd/user/sway-session.target
[Unit]
Description=Sway compositor session
Documentation=man:sytemd.special
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
